# overlay/etc/uci-defaults/98-check-atc-1
#!/bin/sh
set -e

# 98-check-atc-1
# If 99-install-atc exists and luci-proto-atc is already installed,
# remove 99-install-atc and 100-check-atc-2 so installation won't run again.

SCRIPT99=/etc/uci-defaults/99-install-atc
SCRIPT100=/etc/uci-defaults/100-check-atc-2

[ -f "$SCRIPT99" ] || exit 0

if apk info -e luci-proto-atc >/dev/null 2>&1; then
    rm -f "$SCRIPT99" "$SCRIPT100" >/dev/null 2>&1
    logger -t 98-check-atc "luci-proto-atc present; removed $SCRIPT99 and $SCRIPT100"
fi

exit 0
