# overlay/etc/uci-defaults/99-install-atc
#!/bin/sh
set -e

# Minimal 99-install-atc: install all bundled .apk files in /root/atc in a single apk add call

DIR=/root/atc
[ -d "$DIR" ] || exit 0

set --
for pkg in "$DIR"/*.apk; do
    [ -f "$pkg" ] || continue
    set -- "$@" "$pkg"
done

[ "$#" -gt 0 ] || exit 0

# Install all found package files in one apk invocation. Reboot/cleanup handled by 98/100.
apk add --allow-untrusted --no-network "$@"

exit 0