# overlay/etc/uci-defaults/99-install-atc
#!/bin/sh
set -e

PKG_DIR="/root/atc"
REBOOT_FLAG="/tmp/atc-reboot.flag"

# Ensure packages exist
if ! ls "$PKG_DIR"/*.apk >/dev/null 2>&1; then
    logger -t uci-defaults "No .apk packages found in $PKG_DIR"
    exit 1
fi

# Install packages
for pkg in "$PKG_DIR"/*.apk; do
    logger -t uci-defaults "Installing $pkg"
    apk add --allow-untrusted --no-network "$pkg"
done

# Schedule a delayed reboot via helper service (executed outside this script)
touch "$REBOOT_FLAG" 2>/dev/null || logger -t uci-defaults "Warning: cannot create $REBOOT_FLAG"
/etc/init.d/atc-delayed-reboot start >/dev/null 2>&1 || logger -t uci-defaults "Warning: cannot start atc-delayed-reboot"

exit 0