# overlay/etc/uci-defaults/99-install-atc
#!/bin/sh
set -e

# Install pre-bundled ATC packages from /root/atc only when luci-proto-atc is missing
if ! apk info -e luci-proto-atc >/dev/null 2>&1; then
    DIR=/root/atc
    if [ -d "$DIR" ]; then
        for pkg in "$DIR"/*.apk; do
            [ -f "$pkg" ] || continue
            apk add --allow-untrusted "$pkg"
        done
    fi
fi

exit 0
