# overlay/etc/crontabs/root
* * * * * [ -f /tmp/atc-reboot.flag ] && rm -f /tmp/atc-reboot.flag && reboot
